# deploy/hermes-brief-web.service
[Unit]
Description=Hermes Garmin morning brief static web server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=@PROJECT_ROOT@
Environment=PYTHONUNBUFFERED=1
ExecStart=@PROJECT_ROOT@/.venv/bin/python @PROJECT_ROOT@/scripts/serve_brief.py --host 0.0.0.0 --port @PORT@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
